# 4x-ui-blockchain-auth.service
[Unit]
Description=4x-ui Blockchain Auth Microservice
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/4x-ui/blockchain-auth
ExecStart=/opt/4x-ui/blockchain-auth/target/release/blockchain-auth
Restart=always
RestartSec=5
User=root
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 